# etc/udev/rules.d/50-audio-proxy.rules
ACTIO=="add" SUBSYSTEM="sound" ... RUN+="/home/debinan/bin/audio-proxy restart"
ACTIO=="remove" SUBSYSTEM="sound" ... RUN+="/home/debinan/bin/audio-proxy stop"
